# fms.service
[Unit]
Description=FMS Sever Service

[Service]
Type=simple
Restart=always
ExecStart=/home/fms/Scrimmage-FMS/fms_scrimmage/cheesy-arena-lite
WorkingDirectory=/home/fms/Scrimmage-FMS/fms_scrimmage

[Install]
WantedBy=multi-user.target
